# apps/crypto-on-the-hour/deploy/crypto-on-the-hour-worker.service.example
# Crypto On the Hour — BullMQ worker (hourly pipeline + webhooks).
#
# Install:
#   sudo cp deploy/crypto-on-the-hour-worker.service.example /etc/systemd/system/crypto-on-the-hour-worker.service
#   sudo nano /etc/systemd/system/crypto-on-the-hour-worker.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now crypto-on-the-hour-worker
#
# StartLimit* belongs in [Unit] (systemd warns if placed under [Service] on some versions).

[Unit]
Description=Block70 Crypto On the Hour worker (BullMQ hourly pipeline)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=200
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/block70/apps/crypto-on-the-hour
EnvironmentFile=/root/block70/apps/crypto-on-the-hour/.env
ExecStart=/usr/bin/node dist/workers/index.js
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
